# fleet-server/packaging/io.github.pfedotovsky.fleet-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd agent for fleet-server on macOS.

  cp io.github.pfedotovsky.fleet-server.plist ~/Library/LaunchAgents/
  launchctl load -w ~/Library/LaunchAgents/io.github.pfedotovsky.fleet-server.plist

  HOST is intentionally unset by default: fleet-server binds :: and falls back
  to 0.0.0.0 if IPv6 is unavailable. Add HOST only to pin the bind address.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.pfedotovsky.fleet-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/fleet-server</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SERVER_PORT</key>
    <string>3011</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/fleet-server.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/fleet-server.log</string>
</dict>
</plist>
